# deploy/jebbysays.service
[Unit]
Description=jebbysays MCP server
After=network.target

[Service]
Type=simple
User=jebbysays
Group=jebbysays
ExecStart=/var/lib/jebbysays/jebbysays serve --path /var/lib/jebbysays/jebbysays.sqlite3
WorkingDirectory=/var/lib/jebbysays
Environment="MCP_SERVER_URL=https://jebbysays.dev"
Environment="OAUTH_ISSUER=https://clerk.jebbysays.dev"
Environment="OAUTH_CLIENT_ID=ghC0HXqA3jrd3eXj"
Environment="LEPTOS_SITE_ROOT=/var/lib/jebbysays/site"
Environment="LEPTOS_ENV=PROD"
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/jebbysays/jebbysays.log
StandardError=append:/var/log/jebbysays/jebbysays.log

[Install]
WantedBy=multi-user.target
